# server/config/supervisord.d/sensu_api.ini
[program:sensu_api]
command=/usr/bin/sensu-api -c /etc/sensu/config.json -d /etc/sensu/conf.d,/etc/sensu/checks,/etc/sensu/handlers -l /var/log/sensu/sensu-api.log -L info
redirect_stderr=true
