# client/config/supervisord.d/sensu_client.ini
[program:sensu_client]
command=/usr/bin/sensu-client -c /etc/sensu/config.json -d /etc/sensu/conf.d,/etc/sensu/handlers -l /var/log/sensu/sensu-client.log -L info
redirect_stderr=true
