# server/config/supervisor.d/sensu_server.ini
[program:sensu_server]
command=/usr/bin/sensu-server -c /etc/sensu/config.json -d "/etc/sensu/conf.d,/etc/sensu/checks,/etc/sensu/handlers" -l /var/log/sensu/server.log -L info
redirect_stderr=true
