# rabbitmq/config/supervisor.d/rabbitmq.ini
[program:rabbitmq]
command=/opt/rabbitmq_server/sbin/rabbitmq-server
environment=RABBITMQ_EBIN_ROOT="/opt/rabbitmq_server/ebin",RABBITMQ_MNESIA_DIR="/var/lib/rabbitmq/mnesia/rabbit@rabbitmq",RABBITMQ_LOGS="/var/log/rabbitmq/rabbit@server.log",RABBITMQ_SASL_LOGS="/var/log/rabbitmq/rabbit@server-sasl.log",RABBITMQ_CONFIG_FILE="/etc/rabbitmq/rabbitmq",RABBITMQ_ENABLED_PLUGINS_FILE="/etc/rabbitmq/enabled_plugins"
redirect_stderr=true
